# workspace/crontab/laradock
* * * * * laradock /usr/bin/php /var/www/rotaradar-web/artisan schedule:run >> /dev/null 2>&1
* * * * * laradock /usr/bin/php /var/www/trivia-web/artisan schedule:run >> /dev/null 2>&1
